# roles/centos-base/templates/tftp.service.j2
[Unit]
Description=Tftp Server
Requires=tftp.socket
Documentation=man:in.tftpd

[Service]
ExecStart=/usr/sbin/in.tftpd -c -s {{ vm_directory }}
StandardInput=socket

[Install]
Also=tftp.socket
